# roles/deadd/files/deadd.service
[Unit]
Description=Deadd notification center

[Service]
ExecStart=/usr/bin/deadd-notification-center

[Install]
WantedBy=default.target
